# debian/cron.d
0 * * * *    www-data          if test -x /usr/share/drupal8/core/scripts/cron.sh ; then /usr/share/drupal8/core/scripts/cron.sh ; fi
